# config/systemd/user/flexget-forget-seen-rss.timer
[Unit]
Description=Forget flexget logged seen entries for rss feeds

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=default.target
