# config/systemd/user/flexget-database-maintenance.timer
[Unit]
Description=Flexget database maintenance

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
